# etc/services.d/gc/run
#!/usr/bin/with-contenv bash
echo "Starting garbage collector"
exec /docker-tc/bin/gc.sh